# nfs-gssd.service
[Unit]
Description=Secure NFS
Requires=var-lib-nfs-rpc_pipefs.mount
After=syslog.target var-lib-nfs-rpc_pipefs.mount nfs-idmapd.service

[Service]
Type=forking
StandardError=syslog+console
EnvironmentFile=-/etc/sysconfig/nfsfs
ExecStart=/usr/sbin/rpc.gssd $RPCGSSOPTIONS

[Install]
WantedBy=multi-user.target
